# hosts/ignitelinuxsrv1/crontab
* * * * * /home/igniteuser/get_events.sh
